# com.agent-reflection.twitter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-reflection.twitter</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd /Users/nicobailon/Documents/development/agent-reflection &amp;&amp; source ~/.zshrc &amp;&amp; npx tsx ingestion/scripts/twitter-ingest-sqlite.ts >> ~/data/agent-reflection/twitter/ingest.log 2>&amp;1</string>
    </array>
    
    <key>WatchPaths</key>
    <array>
        <string>/Users/nicobailon/data/agent-reflection/twitter</string>
        <string>/Users/nicobailon/data/agent-reflection/twitter/likes</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>10</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Users/nicobailon/.local/bin:/opt/homebrew/bin</string>
        <key>DB_PATH</key>
        <string>/Users/nicobailon/data/agent-reflection/reflection.db</string>
        <key>OPENAI_API_KEY</key>
        <string>YOUR_KEY</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>StandardOutPath</key>
    <string>/Users/nicobailon/data/agent-reflection/twitter/ingest.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/nicobailon/data/agent-reflection/twitter/ingest-error.log</string>
    
    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
